# Scripts/server/DEP/mnt-fastpool.automount
[Unit]
Description=Automount for fastpool
DefaultDependencies=no
Requires=network-online.target
After=network-online.target
Before=remote-fs.target

[Automount]
Where=/mnt/fastpool
TimeoutIdleSec=600

[Install]
WantedBy=multi-user.target
